# br-ext/overlay/etc/init.d/S05calibrate
#!/bin/sh
### BEGIN INIT INFO
# Provides:          touch-calibrate-firstboot
# Required-Start:    $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Short-Description: Run touchscreen calibration on first boot
### END INIT INFO

CAL_MARKER="/etc/.touch_calibrated"
TS_CAL="/data/tslib/ts.cal"
TS_CONF="/data/tslib/ts.conf"
TSPROGRAM="/usr/bin/ts_calibrate"

# Show splash before calibration
if command -v fbi >/dev/null 2>&1; then
    fbi -T 1 -noverbose /etc/splash.ppm >/dev/null 2>&1 &
    sleep 2
fi

# Exit if already calibrated
[ -f "$CAL_MARKER" ] && exit 0

# Exit if ts_calibrate not installed
[ ! -x "$TSPROGRAM" ] && exit 0

# Export env so ts_calibrate finds device and config
export TSLIB_TSDEVICE=/dev/input/touchscreen0
export TSLIB_CALIBFILE=$TS_CAL
export TSLIB_CONFFILE=$TS_CONF

# Run calibration
if command -v start-stop-daemon >/dev/null 2>&1; then
    start-stop-daemon -S -c trezor:trezor --exec "$TSPROGRAM" -- -f $TS_CAL
else
    $TSPROGRAM -f $TS_CAL
fi
RC=$?

# Mark as done
if [ "$RC" -eq 0 ]; then
    touch "$CAL_MARKER"
fi

exit 0
